# rigd_tbot/systemd_units/phase_supervisor.service
[Unit]
[Unit]
Description=RIGD TradeBot Phase Supervisor
After=network.target

[Service]
Type=simple
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
WorkingDirectory=%h/rigd_tbot
ExecStart=/usr/bin/env python3 tbot_bot/support/phase_supervisor.py
Restart=always

[Install]
WantedBy=default.target
