# rigd_tbot/systemd_units/tbot_bot.path
# tbot_bot.path

[Unit]
Description=RIGD TradeBot Core Path Monitor

[Path]
PathChanged=%h/rigd_tbot/tbot_bot/control/control_start.txt

[Install]
WantedBy=default.target
