# rigd_tbot/systemd_units/tbot_web.service
# systemd_units/tbot_web.service
# Launches the Flask web UI as a standalone service
[Unit]
Description=RIGD TradeBot Web UI (Flask)
After=network.target

[Service]
Type=simple
User=tbot
WorkingDirectory=/home/tbot/rigd_tbot
ExecStartPre=/bin/echo [systemd:tbot_web] Launching Flask web UI...
ExecStart=/home/tbot/rigd_tbot/venv/bin/python3 tbot_web/py/run_web.py
ExecStartPost=/bin/echo [systemd:tbot_web] Flask web UI started.
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
