# rigd_tbot/systemd_units/tbot_bot.service
[Unit]
Description=RIGD TradeBot (single-entry: web + scheduler)

[Service]
Type=simple
WorkingDirectory=/home/tbot/rigd_tbot
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -m tbot_bot.runtime.main
Restart=always
RestartSec=5
# (Optional) hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
