# rigd_tbot/systemd_units/tbot_web_bootstrap.service
# tbot_web_bootstrap.service
[Unit]
Description=RIGD TradeBot Web UI (Bootstrap/Config Phase)
After=network.target

[Service]
Type=simple
User=tbot
WorkingDirectory=/home/tbot/rigd_tbot/
ExecStart=/home/tbot/rigd_tbot/venv/bin/python3 -m tbot_web.py.portal_web_bootstrap
Restart=no
Environment="FLASK_SECRET_KEY=changeme-unsafe-dev-key"
Environment="BOT_WEB_PHASE=bootstrap"

[Install]
WantedBy=tbot_web.target
