# rigd_tbot/systemd_units/tbot_web_bootstrap.service
# tbot_web_bootstrap.service

[Unit]
Description=RIGD TradeBot Web UI (Bootstrap/Config Phase)
After=network.target

[Service]
Type=simple
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment="FLASK_SECRET_KEY=changeme-unsafe-dev-key"
Environment="BOT_WEB_PHASE=bootstrap"
WorkingDirectory=%h/rigd_tbot
ExecStart=/usr/bin/env python3 tbot_web/py/portal_web_bootstrap.py
Restart=no

[Install]
WantedBy=default.target
