# rigd_tbot/systemd_units/tbot_provisioning.service
# tbot_provisioning.service

[Unit]
Description=RIGD TradeBot Provisioning Orchestrator (Keys/Secrets Generation Phase)
After=network.target

[Service]
Type=simple
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
WorkingDirectory=%h/rigd_tbot
ExecStart=/usr/bin/env python3 tbot_bot/config/provisioning_runner.py

Restart=no
Environment="PROVISIONING_PHASE=1"

[Install]
WantedBy=default.target
