# rigd_tbot/systemd_units/tbot_provisioning.service
# systemd_units/tbot_provisioning.service
# Launches the provisioning_runner for bootstrap/setup
[Unit]
Description=RIGD TradeBot Provisioning Runner
After=network.target

[Service]
Type=simple
User=tbot
WorkingDirectory=/home/tbot/rigd_tbot
Environment=PYTHONPATH=/home/tbot/rigd_tbot
ExecStartPre=/bin/echo "[systemd:tbot_provisioning] Starting provisioning runner..."
ExecStart=/home/tbot/rigd_tbot/venv/bin/python3 tbot_bot/config/provisioning_runner.py
ExecStartPost=/bin/echo "[systemd:tbot_provisioning] Provisioning runner launched."
Restart=no
RestartSec=2

[Install]
WantedBy=multi-user.target
