# rigd_tbot/systemd_units/tbot_bot.service
# tbot_bot.service
[Unit]
Description=RIGD TradeBot Core Bot Runtime
After=network.target

[Service]
Type=simple
User=tbot
WorkingDirectory=/home/tbot/rigd_tbot/
ExecStart=/home/tbot/rigd_tbot/venv/bin/python3 -m tbot_bot.runtime.main
Restart=on-failure

[Install]
WantedBy=multi-user.target
