# rigd_tbot/systemd_units/tbot_supervisor@.service
[Unit]
Description=RIGD TradeBot Daily Supervisor (%i)
After=network-online.target tbot_bot.service
Wants=network-online.target tbot_bot.service
# block until you’ve configured/provisioned
ConditionPathExists=!/home/tbot/rigd_tbot/tbot_bot/control/provisioning.lock
ConditionPathExists=!/home/tbot/rigd_tbot/tbot_bot/control/bootstrapping.lock

[Service]
Type=oneshot
WorkingDirectory=/home/tbot/rigd_tbot
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -m tbot_bot.runtime.tbot_supervisor --date=%i

# Hardening (kept as-is)
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
UMask=0077

[Install]
WantedBy=default.target
