# rigd_tbot/systemd_units/tbot_supervisor.timer
# systemd_units/tbot_supervisor.timer
[Unit]
Description=RIGD TradeBot Daily Supervisor Timer

[Timer]
# Fire every weekday at 00:01 UTC (explicit form with date + seconds for robustness)
OnCalendar=Mon..Fri *-*-* 00:01:00 UTC
Unit=tbot_supervisor@%Y-%m-%d.service
Persistent=true
AccuracySec=1min
RandomizedDelaySec=0

[Install]
WantedBy=timers.target
