# rigd_tbot/systemd_units/tbot_provisioning.service
# tbot_provisioning.service

[Unit]
Description=RIGD TradeBot Provisioning Orchestrator
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/tbot/rigd_tbot/
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/tbot/rigd_tbot/venv/bin/python3 -m tbot_bot.config.provisioning_runner
Restart=no

[Install]
WantedBy=default.target
