# rigd_tbot/systemd_units/tbot_bot.path
# tbot_bot.path
[Unit]
Description=RIGD TradeBot Core Path Monitor
Wants=tbot_bot.service

[Path]
PathChanged=/home/tbot/rigd_tbot/tbot_bot/control/control_start.txt

[Install]
WantedBy=multi-user.target
