# rigd_tbot/systemd_units/tbot_bot.service
[Unit]
Description=RIGD Tbot Core Bot Runtime (Supervisor Orchestration, v045+)
After=network.target

[Service]
Type=simple
WorkingDirectory=%E/rigd_tbot
Environment=PATH=%E/rigd_tbot/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=%E/rigd_tbot
ExecStart=%E/rigd_tbot/venv/bin/python3 %E/rigd_tbot/tbot_bot/runtime/main.py
Restart=always
RestartSec=2s

# All persistent worker/watcher/test runner modules are orchestrated ONLY by tbot_supervisor.py.
# No process in tbot_bot/ may be launched by CLI, systemd, or main.py except tbot_supervisor.py (v045+).

[Install]
WantedBy=default.target
