# rigd_tbot/systemd_units/tbot_bot.service
[Unit]
Description=RIGD TradeBot Web UI (no scheduler)
After=network.target

[Service]
Type=simple
WorkingDirectory=%E/rigd_tbot
Environment=PYTHONPATH=%E/rigd_tbot
ExecStart=/usr/bin/python3 -m tbot_bot.runtime.main
Restart=always
RestartSec=2s

[Install]
WantedBy=default.target
