# rigd_tbot/systemd_units/tbot_supervisor@.service
# systemd_units/tbot_supervisor@.service
[Unit]
Description=RIGD TradeBot Daily Supervisor (%i) — one-shot orchestration
After=network-online.target tbot_bot.service
Wants=network-online.target tbot_bot.service
# Do NOT run while provisioning/bootstrap locks exist
ConditionPathExists=!%h/rigd_tbot/tbot_bot/control/provisioning.lock
ConditionPathExists=!%h/rigd_tbot/tbot_bot/control/bootstrapping.lock

[Service]
Type=oneshot
WorkingDirectory=%h/rigd_tbot
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=%h/rigd_tbot
ExecStart=/usr/bin/python3 -m tbot_bot.runtime.tbot_supervisor --date=%i

# Hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
UMask=0077

[Install]
# (Template unit is started by the timer; no need to enable it directly)
WantedBy=default.target
